# deploy/systemd/infradocs-fed-reap.service
[Unit]
# PRIMARY ONLY. Sweeps the federation_commands queue: any command a secondary
# claimed (status=dispatched) but never reported back, once older than
# COMMAND_EXPIRY_SECONDS (900s), is closed out to status=expired with a closing
# actions_log entry. Never deletes — the audit trail is preserved. Operates on
# the primary's MongoDB directly; it does not talk to any secondary.
# Self-guarding: `federation_agent reap` exits 1 without acting unless this host's
# settings.role == "primary", so installing it on the wrong host is a harmless
# no-op. Role is read from per-host config, never hardcoded here.
Description=InfraDocs federation reap (primary: expire stale dispatched commands)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
# --- PATHS: adjust these three to match the primary host ----------------------
#   User / WorkingDirectory / the venv interpreter in ExecStart must point at the
#   InfraDocs checkout on the PRIMARY. Values below match the reference deploy
#   (msinha @ ~/projects/InfraDocs_V6). EnvironmentFile supplies INFRADOCS_MONGO_URI.
User=msinha
Group=msinha
WorkingDirectory=/home/msinha/projects/InfraDocs_V6
EnvironmentFile=/home/msinha/projects/InfraDocs_V6/.env
ExecStart=/home/msinha/projects/InfraDocs_V6/venv/bin/python -m app.federation_agent reap
