# deploy/systemd/infradocs-fed-poll.timer
[Unit]
Description=Run InfraDocs federation poll every 5 minutes (secondary)

[Timer]
# First poll shortly after boot, then every 5 minutes. This sets the worst-case
# latency between queuing a command on the primary and the secondary running it.
OnBootSec=2min
OnUnitActiveSec=5min
Persistent=true

[Install]
WantedBy=timers.target
